# etc/systemd/ketchup-rsync.path.in
[Unit]
Description=Watch Claude creds dir and rsync to ketchup relay

[Path]
PathChanged=@CREDS_DIR@
Unit=ketchup-rsync.service

[Install]
WantedBy=default.target
